# salt/schub/templates/docker-controller.conf
description "Docker Controller"

start on filesystem and started docker and started consul
stop on runlevel [!2345]

respawn
respawn limit 10 10

pre-start script
  /usr/bin/docker rm -f docker-controller ||true
  /usr/bin/docker pull {{ pillar['schub']['docker-controller']['image'] }}
end script

script
  /usr/bin/docker run \
    --restart=always \
    --name=docker-controller \
    -v /var/run/docker.sock:/var/run/docker.sock \
    {{ pillar['schub']['docker-controller']['image'] }} -collector=docker:///var/run/docker.sock -storage=consul://{{ pillar['docker']['gw'] }}:8500 -cadvisor=http://{{ pillar['docker']['gw'] }}:4040 -interval=2
end script
